# testRabbitMQServer.conf
#!upstart
description "Test Rabbit MQ server"

start on startup
stop on shutdown

respawn

exec /home/michael-anthony-rodriguez/RabbitMQ/testRabbitMQServer.php

post-start script
    PID=`status testRabbitMQServer | egrep -oi '([0-9]+)$' | head -n1`
    echo $PID > /testRabbitMQServer.pid
end script

post-stop script
    rm -f /testRabbitMQServer.pid
end script
