# dev/systemd/qa-backend.service
[Unit]
Description=QA Backend listener
Wants=backend-server.service
After=backend-server.service network.target

[Service]
Type=simple
WorkingDirectory=%h/Cinemaniacs/vm
ExecStart=/usr/bin/php %h/Cinemaniacs/vm/server.php QA backend
Restart=always
RestartSec=5

[Install]
WantedBy=default.target
